# scripts/launchd/com.rva-figure-drawing.discover.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.rva-figure-drawing.discover</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/sean/sites/rva-figure-drawing-calendar/scripts/discover-sources.sh</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/sean/sites/rva-figure-drawing-calendar</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Weekday</key>
		<integer>0</integer>
		<key>Hour</key>
		<integer>9</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/Users/sean/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
	</dict>
	<key>StandardOutPath</key>
	<string>/Users/sean/sites/rva-figure-drawing-calendar/logs/discover-stdout.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/sean/sites/rva-figure-drawing-calendar/logs/discover-stderr.log</string>
</dict>
</plist>
